# com.jobfitevaluator.poller.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd agent: run the poller daily at 8am.
     Setup:
       1. Edit the paths below to match your machine (the repo path and your home dir).
       2. Copy to ~/Library/LaunchAgents/com.jobfitevaluator.poller.plist
       3. launchctl load -w ~/Library/LaunchAgents/com.jobfitevaluator.poller.plist
     Trigger a test run:  launchctl start com.jobfitevaluator.poller
     Remove:              launchctl unload ~/Library/LaunchAgents/com.jobfitevaluator.poller.plist -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jobfitevaluator.poller</string>
    <key>ProgramArguments</key>
    <array>
        <string>/ABSOLUTE/PATH/TO/job-fit-evaluator/.venv/bin/python</string>
        <string>/ABSOLUTE/PATH/TO/job-fit-evaluator/poller.py</string>
        <string>--bay-area-only</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/ABSOLUTE/PATH/TO/job-fit-evaluator</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/YOU/.job-fit-evaluator/poller.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOU/.job-fit-evaluator/poller.log</string>
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
